# gunicorn.service
[Unit]
Description=Gunicorn fucks Django /etc/systemd/system/gunicorn.service
Before=nginx.service
After=network.target

[Service]
WorkingDirectory=/home/django/deploy/app
ExecStart=/usr/bin/gunicorn3 --name=app --pythonpath=/home/django/deploy/app/venv/lib/python3.8/site-packages --bind unix:/home/django/gunicorn.socket --config /etc/gunicorn.d/gunicorn.py app.wsgi:application
Restart=always
SyslogIdentifier=gunicorn
User=django
Group=django


[Install]
WantedBy=multi-user.target
